# prime-database-backup/backup.cron
# Runs backup 
59 1 * * * sh /opt/backup/backup.sh 
