# prime-database-backup/backup.cron
# Runs backup 
*/1 * * * * sh /opt/backup/backup.sh 
